# docker/supervisor/worker/horizon.conf
[program:horizon]
command=php /var/www/artisan horizon
directory=/var/www
autostart=true
autorestart=true
user=www-data
stdout_logfile=/var/log/supervisor/horizon.log
redirect_stderr=true
